# rootfs/etc/s6-overlay/s6-rc.d/svc-kclient/run
#!/usr/bin/with-contenv with-user

# NodeJS wrapper
cd /kclient
exec node index.js
