# rootfs/etc/s6-overlay/s6-rc.d/svc-kclient/run
#!/usr/bin/with-contenv with-user

if [[ "${SVC_KCLIENT}" == "0" ]]; then
    echo 'disable kclient serve ...'
    sleep 1
    sudo s6-rc stop svc-kclient
    exit
fi

# DBus Setup
sudo dbus-daemon --system

# Mic Setup
# if [ ! -f '/dev/shm/mic.lock' ]; then
#   until [ -f /defaults/pid ]; do
#     sleep .5
#   done
#   s6-setuidgid ${USER} with-contenv pactl \
#     load-module module-pipe-source \
#     source_name=virtmic \
#     file=/defaults/mic.sock \
#     source_properties=device.description=LSIOMic \
#     format=s16le \
#     rate=44100 \
#     channels=1
#   s6-setuidgid ${USER} with-contenv pactl \
#     set-default-source virtmic
#   touch /dev/shm/mic.lock
# fi

# NodeJS wrapper
cd /kclient
exec node index.js
