# rootfs/etc/s6-overlay/s6-rc.d/init-services/run
#!/usr/bin/with-contenv bash

# 更改文件权限，否则 firefox 无法启动
if [[ "$USER_DEF" == "root" ]] || [[ "$USER_DEF" == "" ]] ; then
    # chown -R root:root $HOME
    echo "init services server by root"
else
    echo "init services server by user"
    sudo chown -R $USER_DEF:$USER_DEF $HOME
fi
