# rootfs/etc/s6-overlay/s6-rc.d/svc-de/run
#!/usr/bin/with-contenv with-user

cd $HOME
exec bash /defaults/startwm.sh
