# rootfs/etc/s6-overlay/s6-rc.d/svc-pulseaudio/run
#!/usr/bin/with-contenv with-user

if [[ "${SVC_KASM}" == "0" ]]; then
    echo 'disable kams[pulseaudio] server ...'
    s6-svc -O exit
fi


exec /usr/bin/pulseaudio \
    --log-level=0 \
    --log-target=stderr \
    --exit-idle-time=-1 > /dev/null 2>&1
