# rootfs/etc/s6-overlay/s6-rc.d/svc-pulseaudio/run
#!/usr/bin/with-contenv with-user

if [[ "${SVC_AUDIO}" == "0" ]]; then
    echo 'disable pulseaudio server ...'
    sleep 1
    sudo s6-rc stop svc-pulseaudio
    exit
fi


exec /usr/bin/pulseaudio \
    --log-level=0 \
    --log-target=stderr \
    --exit-idle-time=-1 > /dev/null 2>&1
