# root/etc/s6-overlay/s6-rc.d/svc-upload/run
#!/command/with-contenv bash
set -e

echo 'Starting upload server'
$STARTUPDIR/upload_server/kasm_upload_server --ssl --auth-token "kasm_user:$VNC_PW"
