# rootfs/etc/s6-overlay/s6-rc.d/init-conf/run
#!/usr/bin/with-contenv with-user
set -e
echo "Init System By S6-Overlay"

# 更改文件权限，否则 firefox 无法启动
if [ "$USERNAME" == "root" ] || [ "$USERNAME" == "" ] ; then
    chown -R root:root $HOME
else
    sudo chown -R $USERNAME:$USERNAME $HOME
fi


