# rootfs/etc/s6-overlay/s6-rc.d/init-conf/run
#!/usr/bin/with-contenv with-user
set -e
echo "Init System By S6-Overlay"
