# rootfs/etc/s6-overlay/s6-rc.d/svc-kasmde/run
#!/usr/bin/with-contenv with-user

if [[ "${SVC_KASM}" == "0" ]]; then
    echo 'disable kams[de] server ...'
    sleep 1
    sudo s6-rc stop svc-kasmde
fi

dbus-daemon --system
cd $HOME
exec bash /defaults/startwm.sh
