# rootfs/etc/s6-overlay/s6-rc.d/svc-kasmde/run
#!/usr/bin/with-contenv with-user

if [[ "${SVC_KASM}" == "0" ]]; then
    echo 'disable kams[de] server ...'
    sleep 1
    sudo s6-rc stop svc-kasmde
fi

echo "startwm[svc-kamsde] by $USER ============="
cd $HOME
sudo chmod +x  /defaults/startwm.sh
exec /bin/bash /defaults/startwm.sh
