# script/webtop.kas.ini
; ============================================================================
; kas.ini - kas 进程管理器配置
; 使用方式: kas -c /wsc/kway/kas.ini -s /wsc/kway/kas.sock
;
; 服务启动顺序:
;   svc-ainit (一次性)
;     -> svc-dbus
;     -> svc-audio
;     -> svc-xorg
;         -> svc-selkies  (需要 xorg)
;             -> svc-kin  (对外端口 7090)
;         -> svc-desktop (需要 xorg + dbus)
; ============================================================================

[program:init]
command=/usr/srv/init.sh
type=once
priority=1
autostart=true
autorestart=false
; stdout_logfile=/tmp/init.log
; stderr_logfile=/tmp/init.err

[program:sshd]
command=/usr/sbin/sshd -D
type=long
priority=10
autostart=${KAS_SSHD}
autorestart=true
depends=init
; stdout_logfile=/tmp/sshd.log
; stderr_logfile=/tmp/sshd.err

[program:vscode]
command=/usr/srv/vscode.sh
type=long
priority=11
autostart=${KAS_VSCODE}
autorestart=true
; user=${USER}
; stdout_logfile=/tmp/vscode.log
; stderr_logfile=/tmp/vscode.err

; -------- svc-ainit: 一次性初始化 (首次运行) --------
[program:svc-ainit]
command=/usr/srv/svc-ainit.sh
autostart=${KAS_WEBTOP}
autorestart=false
type=once
priority=1
; stdout_logfile=/tmp/svc/ainit.log
; stderr_logfile=/tmp/svc/ainit.err

; -------- svc-dbus: D-Bus 系统守护进程 --------
[program:svc-dbus]
command=/usr/srv/svc-dbus.sh
autostart=${KAS_WEBTOP}
autorestart=true
type=long
priority=5
depends=svc-ainit
; stdout_logfile=/tmp/svc/dbus.log
; stderr_logfile=/tmp/svc/dbus.err

; -------- svc-audio: PulseAudio 音频守护进程 --------
[program:svc-audio]
command=/usr/srv/svc-audio.sh
autostart=${KAS_WEBTOP}
autorestart=true
type=long
priority=6
depends=svc-ainit
; stdout_logfile=/tmp/svc/audio.log
; stderr_logfile=/tmp/svc/audio.err

; -------- svc-xorg: Xvfb 虚拟显示 --------
[program:svc-xorg]
command=/usr/srv/svc-xorg.sh
autostart=${KAS_WEBTOP}
autorestart=true
type=long
priority=10
depends=svc-ainit
; stdout_logfile=/tmp/svc/xorg.log
; stderr_logfile=/tmp/svc/xorg.err

; -------- svc-selkies: Selkies WebRTC 桥接 (端口 8082) --------
[program:svc-selkies]
command=/usr/srv/svc-selkies.sh
autostart=${KAS_WEBTOP}
autorestart=true
type=long
priority=20
depends=svc-xorg,svc-audio
stopwaitsecs=10
max_retries=5
restart_delay=3
; stdout_logfile=/tmp/svc/selkies.log
; stderr_logfile=/tmp/svc/selkies.err

; -------- svc-desktop: XFCE 桌面环境 --------
[program:svc-desktop]
command=/usr/srv/svc-desktop.sh
autostart=${KAS_WEBTOP}
autorestart=true
type=long
priority=20
depends=svc-xorg,svc-dbus
stopwaitsecs=10
max_retries=5
restart_delay=3
; stdout_logfile=/tmp/svc/desktop.log
; stderr_logfile=/tmp/svc/desktop.err

; -------- svc-gateway: Kin 反向代理 (对外端口 7090, 无需密码) --------
[program:svc-gateway]
command=/usr/srv/svc-gateway.sh
autostart=${KAS_WEBTOP}
autorestart=true
type=long
priority=30
depends=svc-selkies
stopwaitsecs=5
max_retries=10
restart_delay=2
; stdout_logfile=/tmp/svc/gateway.log
; stderr_logfile=/tmp/svc/gateway.err
