# rootfs/etc/s6-overlay/s6-rc.d/svc-sshd/run
#!/usr/bin/with-contenv bash

if [[ "${SVC_SSHD}" == "0" ]]; then
    echo 'disable sshd server ...'
    s6-rc vsc-sshd stop
    exit
fi

if [ ! -f '/etc/ssh/ssh_host_rsa_key' ]; then
    echo 'init sshd keygen ...'
    # sshd-keygen,初始化
    ssh-keygen -q -t rsa -N '' -f /etc/ssh/ssh_host_rsa_key
    ssh-keygen -q -t rsa -N '' -f /etc/ssh/ssh_host_dsa_key
    ssh-keygen -q -t rsa -N '' -f /etc/ssh/ssh_host_ecdsa_key
    ssh-keygen -q -t rsa -N '' -f /etc/ssh/ssh_host_ed25519_key
    # 配置root的登录密码
    echo "user:${PASSWORD}" | chpasswd
    echo "PasswordAuthentication yes" >> /etc/ssh/sshd_config
    # echo "PermitRootLogin yes" >> /etc/ssh/sshd_config
    mkdir /run/sshd
fi

echo 'start sshd server ... '
exec /usr/sbin/sshd -D
